# deploy/einsatzleiter.service
[Unit]
Description=Einsatzleiter-Hilfswerkzeug FF Wolfurt
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=exec
User=clp-einsatz
Group=clp-einsatz
WorkingDirectory=/home/clp-einsatz/htdocs/einsatzleiter
EnvironmentFile=/home/clp-einsatz/htdocs/einsatzleiter/.env
ExecStart=/home/clp-einsatz/htdocs/einsatzleiter/.venv/bin/gunicorn \
    -k uvicorn.workers.UvicornWorker \
    -w 2 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --keepalive 5 \
    --access-logfile - \
    --error-logfile - \
    app.main:app
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=einsatzleiter

[Install]
WantedBy=multi-user.target
